# MonitoringAgent/packaging/deb/DEBIAN/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable monitoring_agent.service
systemctl start monitoring_agent.service || :
